# roles/linroom/files/99-wakeonlan.rules
ACTION=="add", SUBSYSTEM=="net", NAME=="e*", RUN+="/usr/sbin/ethtool -s $name wol g"
